# config/backupmynode.timer
[Unit]
Description=Run BackupMynode Every 15 minutes
[Timer]
OnBootSec=300
OnUnitActiveSec=86400
AccuracySec=1
[Install]
WantedBy=timers.target
